# crontab
30 22 * * * root /usr/local/bin/scrape.sh
